# config/supervisord.conf
[supervisord]
nodaemon=true

[program:elasticsearch]
directory=/data
command=/elasticsearch/bin/elasticsearch -Des.config=/etc/elasticsearch.yml
autorestart = true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log